# actividad3/saludo.service
[Unit]
Description=Saludo y fecha actual infinitos

[Service]
Type=simple
ExecStart=
Restart=always

[Install]
WantedBy=default.target
